# templates/rapidpen-log-cleanup.timer.template
[Unit]
Description=Run RapidPen Edge Log Cleanup daily

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=3600

[Install]
WantedBy=timers.target
