# templates/rapidpen-log-cleanup.service.template
[Unit]
Description=RapidPen Edge Log Cleanup

[Service]
Type=oneshot
ExecStart=/usr/local/bin/rapidpen-log-cleanup.sh
